# rootfs/etc/systemd/system/rpm-ostreed-automatic.timer.d/override.conf
[Timer]
RandomizedDelaySec=10m
OnCalendar=*-*-* 7:30:00
Persistent=true