# rootfs/usr/etc/systemd/system/rpm-ostreed-automatic.timer.d/override.conf
[Timer]
RandomizedDelaySec=10m
OnBootSec=5m
OnUnitActiveSec=3h
Persistent=true